# rashberrypi/home/user/cron/user.cron

@reboot sleep 500; /home/user/bin/TransferFiles.sh  >> /tmp/TransferFiles.log  2>&1

